# scripts/systemd/browserai-backup.service
[Unit]
Description=BrowserAI — snapshot DeepSeek session file
Wants=browserai-backup.timer

[Service]
Type=oneshot
Environment=DATA_DIR=/opt/browserai-data
Environment=BACKUP_DIR=/var/backups/browserai
Environment=KEEP_DAYS=7
ExecStart=/opt/browserai/scripts/backup-deepseek-session.sh

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/backups/browserai
ProtectHome=true
