# scripts/systemd/browserai-backup.timer
[Unit]
Description=BrowserAI — run SQLite backup daily at 02:30 (Step 10.9)

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
